# code/mountain-beacon.service
[Unit]
Description=NIST IR 8213 Beacon - Mountain RNG
After=network.target mountain-rng.service
Wants=mountain-rng.service

[Service]
Type=simple
User=iip
WorkingDirectory=/home/iip/mountain-rng/beacon
ExecStart=/usr/bin/python3 /home/iip/mountain-rng/beacon/beacon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
